# assets/launch_agents/com.local.doh.plist.j2
<!-- com.local.doh.plist.j2 -->
<plist version="1.0">
<dict>
    <key>Label</key><string>com.local.doh</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{CLOUDFLARED_BIN}}</string>
        <string>proxy-dns</string>
        <string>--address</string><string>127.0.0.1</string>
        <string>--port</string><string>5053</string>
        <string>--upstream</string><string>https://1.1.1.1/dns-query</string>
        <string>--upstream</string><string>https://1.0.0.1/dns-query</string>
    </array>

    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>

    <key>StandardOutPath</key>
    <string>{{LOG_PATH}}</string>
    <key>StandardErrorPath</key>
    <string>{{LOG_PATH}}</string>
</dict>
</plist>
